# launchd/com.example.three-headed-snake-orchestrator.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.example.three-headed-snake-orchestrator</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__REPO_ROOT__/scripts/start-orchestrator.sh</string>
    <string>--poll-ms</string>
    <string>200</string>
  </array>
  <key>WorkingDirectory</key>
  <string>__REPO_ROOT__</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>__REPO_ROOT__/logs/orchestrator.launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>__REPO_ROOT__/logs/orchestrator.launchd.err.log</string>
</dict>
</plist>
